# systemd/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin sysconfig --noclear %I $TERM
